# crontab
# m h dom mon dow   command
0 8 * * * bash /usr/local/bin/backup-files.sh
0 8 * * * bash /usr/local/bin/backup-sql.sh
